# wallpaper-update.timer
[Unit]
Description=Update dual screen wallpaper

[Timer]
OnUnitActiveSec={WALLPAPER_TIMER}
OnBootSec=60s
Unit=wallpaper-update.service
Persistent=true

[Install]
WantedBy=timers.target
